# scripts/systemd/content-processor.service
[Unit]
Description=Reddit Content Processor with LLM Analysis
Documentation=file:///opt/aidaily-digest/scripts/README.md
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=1h

[Service]
Type=oneshot
User=asher
Group=asher
WorkingDirectory=/opt/aidaily-digest/scripts
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run entrypoints/process_content.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=content-processor

# Resource limits (adjusted for smaller VPS)
MemoryMax=500M
CPUQuota=50%

# Restart policy for failures
Restart=on-failure
RestartSec=60

# Timeout settings (LLM calls can take time)
TimeoutStartSec=300
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/aidaily-digest/scripts