# scripts/systemd/reddit-fetcher.service
[Unit]
Description=Reddit AI Content Fetcher
Documentation=file:///opt/aidaily-digest/scripts/README.md
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=10m

[Service]
Type=oneshot
User=asher
Group=asher
WorkingDirectory=/opt/aidaily-digest/scripts
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=UV_CACHE_DIR=/opt/aidaily-digest/scripts/.uv-cache
ExecStart=/usr/local/bin/uv run entrypoints/run_daily_fetch.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reddit-fetcher

# Resource limits (adjusted for smaller VPS)
MemoryMax=300M
CPUQuota=30%

# Restart policy for failures
Restart=on-failure
RestartSec=30

# Security settings
#NoNewPrivileges=yes
#PrivateTmp=yes
#ProtectSystem=yes
#ProtectHome=yes

[Install]
WantedBy=multi-user.target
