# scripts/systemd/twitter-poster.service
[Unit]
Description=Twitter Poster - Automated tweeting of high-scoring Reddit posts
After=network.target

[Service]
Type=oneshot
User=asher
Group=asher
WorkingDirectory=/opt/aidaily-digest/scripts
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=UV_CACHE_DIR=/opt/aidaily-digest/scripts/.uv-cache
ExecStart=/usr/local/bin/uv run entrypoints/run_twitter_poster.py
StandardOutput=journal
StandardError=journal

# Environment file containing secrets
EnvironmentFile=-/opt/aidaily-digest/scripts/.env

# Resource limits (adjusted for smaller VPS)
MemoryMax=300M
CPUQuota=30%

[Install]
WantedBy=multi-user.target